# chrono.service
[Unit]
Description=Chrono Consulting SQL game
After=network.target

[Service]
Type=simple
User=chrono
WorkingDirectory=/opt/chrono/app
EnvironmentFile=/opt/chrono/app/.env
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/chrono

[Install]
WantedBy=multi-user.target
